# scripts/safework-monitor.service
[Unit]
Description=SafeWork 자동 트러블슈팅 모니터링 시스템
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=jclee
Group=jclee
WorkingDirectory=/home/jclee/app/safework2
ExecStart=/home/jclee/app/safework2/scripts/monitor-daemon.sh
Restart=always
RestartSec=300
StandardOutput=journal
StandardError=journal

# 환경변수 설정
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="SLACK_WEBHOOK_URL="
Environment="DISCORD_WEBHOOK_URL="

[Install]
WantedBy=multi-user.target